# wifi_fallback.service
[Unit]
Description=WiFi Fallback Hotspot
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi_fallback.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target