# wifi_fallback.service
[Unit]
Description=WiFi Fallback to Hotspot
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/wifi_fallback.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
