# wifi_fallback.service
[Unit]
Description=WiFi Fallback to Hotspot
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/wifi_fallback.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
