# wifi_fallback.service
[Unit]
Description=WiFi Fallback to Hotspot
After=network-online.target NetworkManager.service
Wants=network-online.target
# Non partire se già connessi
ConditionPathExists=!/tmp/hotspot_active

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/wifi_fallback.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120
KillMode=process

[Install]
WantedBy=multi-user.target